# crontab
* * * * * php /var/www/utils/fetch_recordings.php >> /var/www/logs/fetch_logs
0 2 * * * php /var/www/utils/download_recordings.php >> /var/www/logs/download_logs
# valid line